# sigma/persistence/scheduled_task_creation.yml
title: Scheduled Task Created or Modified (Suspicious Payload)
id: 6f0b29a0-ef58-4c85-9a0a-4f1f1c2d540e
status: stable
description: Detects scheduled task creation/modification via schtasks where the payload points to suspicious locations or LOLBins.
author: N SV
date: 2025-12-20
references:
  - https://attack.mitre.org/techniques/T1053/005/
tags:
  - attack.persistence
  - attack.t1053.005
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    Image|endswith: '\schtasks.exe'
    CommandLine|contains:
      - '/create'
      - '/change'
  selection_suspicious:
    CommandLine|contains:
      - '\AppData\'
      - '\Temp\'
      - 'powershell'
      - 'cmd.exe /c'
      - 'rundll32'
  condition: selection and selection_suspicious
fields:
  - Image
  - CommandLine
  - User
falsepositives:
  - Legitimate IT automation that uses schtasks (baseline + allowlist task names)
level: medium
